# deploy/audio-router.service
[Unit]
Description=Audio Router Service
After=network.target

[Service]
User=audio
Group=audio
WorkingDirectory=/opt/audio-router
ExecStart=/opt/audio-router/AudioRouter --provider JackSharpCore
Restart=on-failure
Environment=ASPNETCORE_ENVIRONMENT=Production

[Install]
WantedBy=multi-user.target
